# amdgpu_workaround/amdgpu_workaround.service
 [Unit]
Description=workaround for AMD GPU

[Service]
User=root
Group=root
ExecStart=/opt/amdgpu_workaround.sh

[Install]
WantedBy=multi-user.target
